# myproject.conf
start on filesystem and started networking
respawn
chdir /root/nodejs/  #deployment directory
exec node-debug --web-port=8080 --web-host=0.0.0.0 myproject.js
